# .config/systemd/user/autostart.target
[Unit]
Description=Run autostart generator
Requires=default.target
After=default.target
Wants=xdg-desktop-autostart.target
