# .config/systemd/user/bg.service
[Unit]
Description=sway output wallpaper rotate service

[Service]
ExecStart=%h/.config/sway/bg.sh

[Install]
WantedBy=multi-user.target
